# luci-app-sms-manager/root/etc/uci-defaults/off_sms_manager.sh
#!/bin/sh
# 
# Copyright 2026 Rafał Wabik (IceG) - From eko.one.pl forum
# Licensed to the GNU General Public License v3.0.
# 

/etc/init.d/sms_manager stop >/dev/null 2>&1 &
/etc/init.d/sms_manager disable >/dev/null 2>&1 &
/etc/init.d/sms_manager disable >/dev/null 2>&1 &
 
exit 0
